# scheduler/scheduler.cron
# run bash script to collect daily stock prices and NAVs at 4:01 pm, m-f, then run predictions
1 16 * * mon-fri ~/StackReport/scheduler/scheduler.sh >& ~/logfile
